# supervisor/config/rabbitmq-worker.conf
[program:rabbitmq-worker]
command=php /var/www/artisan rabbitmq:consume
autostart=true
autorestart=true
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/rabbitmq-worker.log
